# netrater.service
[Unit]
Description=NetRater – AIMD Bandwidth Controller
Documentation=https://github.com/vonarian/netrater
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStartPre=/usr/local/bin/setup-qos.sh
ExecStart=/usr/local/bin/netrater
ExecStop=/usr/local/bin/teardown-qos.sh
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netrater

[Install]
WantedBy=multi-user.target